# infra/systemd/home-school-management-control-api.service
[Unit]
Description=Home School Management Control API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/debian/apps/home-school-management/control-api
Environment=CONTROL_APP_ENV=production
Environment=CONTROL_APP_PORT=3100
Environment=CONTROL_APP_CORS_ORIGIN=http://192.168.1.210
Environment=CONTROL_SESSION_COOKIE_NAME=hsm_operator_session
Environment=CONTROL_SESSION_COOKIE_SECURE=false
Environment=CONTROL_SESSION_COOKIE_SAMESITE=Lax
Environment=CONTROL_SESSION_TTL_HOURS=12
Environment=CONTROL_INTERNAL_API_KEY=HSM_CONTROL_PLANE_SYNC_KEY
Environment=CONTROL_PGHOST=192.168.1.202
Environment=CONTROL_PGPORT=5432
Environment=CONTROL_PGDATABASE=appdb
Environment=CONTROL_PGUSER=appuser
Environment=CONTROL_PGPASSWORD=StrongPassword123!
Environment="PGOPTIONS=-c search_path=hsm_control_staging"
ExecStart=/usr/bin/node src/app.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
